# docker/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:my-program]
command=/usr/local/bin/supervisord
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisord.log
